# deploy/ctx-backup.service
[Unit]
Description=ctx — daily SQLite backup (the store is the sole corpus)
After=network-online.target

[Service]
Type=oneshot
User=synxadmin
Group=synxadmin
WorkingDirectory=/home/synxadmin/Exposemd
Environment=CTX_DATA_DIR=/home/synxadmin/.ctx-data
ExecStart=/home/synxadmin/Exposemd/.venv/bin/ctx backup
