# deploy/ctx.service
[Unit]
Description=ctx — token-efficient .md routing (MCP + UI)
Documentation=https://ctx.prod.synergix.ch
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=synxadmin
Group=synxadmin
WorkingDirectory=/home/synxadmin/Exposemd
EnvironmentFile=-/home/synxadmin/.ctx-data/env
Environment=PATH=/home/synxadmin/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=CTX_PROJECT_ROOT=/home/synxadmin/synergix_prod
Environment=CTX_DATA_DIR=/home/synxadmin/.ctx-data
Environment=CTX_HOST=0.0.0.0
Environment=CTX_PORT=8770
Environment=FASTEMBED_CACHE_PATH=/home/synxadmin/.ctx-data/fastembed_cache
Environment=HF_HOME=/home/synxadmin/.ctx-data/hf_cache
ExecStart=/home/synxadmin/Exposemd/.venv/bin/uvicorn ctx.server:build_app --factory --host 0.0.0.0 --port 8770
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/synxadmin/.ctx-data /home/synxadmin/.cache

[Install]
WantedBy=multi-user.target
