# deploy/trovex-backup.service
[Unit]
Description=trovex — daily SQLite backup (the store is the sole corpus)
After=network-online.target

[Service]
Type=oneshot
User=trovex
Group=trovex
WorkingDirectory=/home/trovex/trovex
Environment=TROVEX_DATA_DIR=/home/trovex/.trovex-data
ExecStart=/home/trovex/trovex/.venv/bin/trovex backup
