# deploy/trovex.service
[Unit]
Description=trovex — token-efficient .md routing (MCP + UI)
Documentation=https://trovex.example.com
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=trovex
Group=trovex
WorkingDirectory=/home/trovex/trovex
EnvironmentFile=-/home/trovex/.trovex-data/env
Environment=PATH=/home/trovex/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=TROVEX_PROJECT_ROOT=/home/trovex/docs-repo
Environment=TROVEX_DATA_DIR=/home/trovex/.trovex-data
Environment=TROVEX_HOST=0.0.0.0
Environment=TROVEX_PORT=8770
Environment=FASTEMBED_CACHE_PATH=/home/trovex/.trovex-data/fastembed_cache
Environment=HF_HOME=/home/trovex/.trovex-data/hf_cache
ExecStart=/home/trovex/trovex/.venv/bin/uvicorn trovex.server:build_app --factory --host 0.0.0.0 --port 8770
Restart=on-failure
RestartSec=5

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/home/trovex/.trovex-data /home/trovex/.cache

[Install]
WantedBy=multi-user.target
