# deploy/trovex-reindex.service
[Unit]
Description=trovex reindex — pull docs-repo and reindex .md docs
After=network-online.target

[Service]
Type=oneshot
User=trovex
Group=trovex
WorkingDirectory=/home/trovex/trovex
EnvironmentFile=-/home/trovex/.trovex-data/env
Environment=PATH=/home/trovex/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=TROVEX_DATA_DIR=/home/trovex/.trovex-data
Environment=FASTEMBED_CACHE_PATH=/home/trovex/.trovex-data/fastembed_cache
Environment=HF_HOME=/home/trovex/.trovex-data/hf_cache
ExecStartPre=-/usr/bin/git -C /home/trovex/docs-repo fetch --quiet
ExecStartPre=-/usr/bin/git -C /home/trovex/docs-repo reset --hard origin/main
ExecStartPre=-/usr/bin/git -C /home/trovex/obsidian/docs-repo pull --quiet --rebase
ExecStartPre=-/usr/bin/git -C /home/trovex/obsidian/docs-repo-dev pull --quiet --rebase
ExecStart=/home/trovex/trovex/.venv/bin/trovex index
